# security/crowdsec-blocklist-mirror/files/crowdsec_mirror.in
#!/bin/sh
#
# PROVIDE: crowdsec_mirror
# REQUIRE: LOGIN DAEMON NETWORKING
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# crowdsec_mirror_enable (bool):	Set it to YES to enable the blocklist mirror.
#					Default is "NO"
# crowdsec_mirror_config (str):		Set the config path.
#					Default is "%%ETCDIR%%/crowdsec-blocklist-mirror.yaml"
# crowdsec_mirror_flags (str):		extra flags to run bouncer.
#					Default is ""

. /etc/rc.subr

name=crowdsec_mirror
desc="Crowdsec Blocklist Mirror"
rcvar=crowdsec_mirror_enable

load_rc_config $name

: "${crowdsec_mirror_enable:=NO}"
: "${crowdsec_mirror_config:=%%ETCDIR%%/crowdsec-blocklist-mirror.yaml}"
: "${crowdsec_mirror_flags:=}"

pidfile=/var/run/${name}.pid
required_files="$crowdsec_mirror_config"
command="%%PREFIX%%/bin/crowdsec-blocklist-mirror"
start_cmd="${name}_start"
start_precmd="${name}_precmd"

crowdsec_mirror_precmd() {
    CSCLI=%%PREFIX%%/bin/cscli
    orig_line="api_key: \${API_KEY}"
    # IF the bouncer is not configured
    if grep -q "^${orig_line}" "${crowdsec_mirror_config}"; then
        SUFFIX=$(LC_CTYPE=C tr -dc A-Za-z0-9 </dev/urandom | head -c 8)
        BOUNCER="cs-blocklist-mirror-${SUFFIX}"
        # AND crowdsec is installed..
        if command -v "$CSCLI" >/dev/null; then
            # THEN, register it to the local API
            API_KEY=$($CSCLI bouncers add "${BOUNCER}" -o raw)
            if [ -n "$API_KEY" ]; then
                sed -i "" "s/^${orig_line}/api_key: ${API_KEY}     # ${BOUNCER}/" "${crowdsec_mirror_config}"
                sed -i "" "s/^${orig_line}/lapi_url: ${CROWDSEC_LAPI_URL}" "http://127.0.0.1:8080"
                echo "Registered: ${BOUNCER}"
            fi
        fi
    fi
}

crowdsec_mirror_start() {
    /usr/sbin/daemon -f -p ${pidfile} -t "${desc}" -- \
        ${command} -c "${crowdsec_mirror_config}" ${crowdsec_mirror_flags}
}

run_rc_command "$1"
